# spec/fixtures/integration/provider/cron/crontab/moved_cronjob_input2
# HEADER: some simple
# HEADER: header
# Oregano Name: some_unrelevant job
* * * * * /bin/true
# Oregano Name: My daily failure
@daily /bin/false
