# spec/fixtures/integration/provider/cron/crontab/purged
# HEADER: some simple
# HEADER: header

# commend with blankline above and below


# Oregano Name: only managed entry
* * * * * /bin/true
